# tor-bw.conf
[program:tor-bw]
command=/usr/bin/python /srv/tor-bw.py
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/tor-bw.log
stderr_logfile=/var/log/tor-bw.log
